# tasmota-master/rootfs/etc/services.d/tasmota-master/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Tasmota Master service stopped"