# tasmota-master/rootfs/etc/services.d/tasmota-master/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Tasmota Master..."

# Wait for dependencies
bashio::log.info "Waiting for dependencies..."
sleep 5

# Change to app directory
cd /opt/app/backend

# Set environment variables
export PYTHONPATH=/opt/app/backend:$PYTHONPATH
export FLASK_ENV=production

# Start the application
bashio::log.info "Starting Flask application..."
exec python3 app.py