# tasmota-master/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting nginx..."

# Start nginx
exec nginx -g "daemon off;"